# etc/supervisord.d/openhab.ini
[program:openhab]
command = /opt/qnib/openhab/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
